# ops/systemd/gta-rp-postgres-backup.timer
[Unit]
Description=Run GTA RP PostgreSQL backup nightly

[Timer]
OnCalendar=*-*-* 03:15:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
