# ops/systemd/gta-rp-postgres-backup.service
[Unit]
Description=GTA RP PostgreSQL backup
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=codex-deploy
Group=codex-deploy
WorkingDirectory=/var/www/gta-rp-population-graph/current
ExecStart=/var/www/gta-rp-population-graph/current/scripts/backup-postgres.sh
