# ops/systemd/gta-rp-uploads-backup.service
[Unit]
Description=GTA RP uploads backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=codex-deploy
Group=codex-deploy
WorkingDirectory=/var/www/gta-rp-population-graph/current
ExecStart=/var/www/gta-rp-population-graph/current/scripts/backup-uploads.sh
